# conf/systemd/user/pulseaudio.service.d/verboselog.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pulseaudio -v --daemonize=no --log-target=journal
