# etc/systemd/system/geoclue.service.d/debug.conf
[Service]
Environment=G_MESSAGES_DEBUG=all
