# etc/systemd/system/earlyoom.service.d/conflicts-oomd.conf
[Unit]
Conflicts=systemd-oomd.service systemd-oomd.socket
