# etc/systemd/system/earlyoom.service.d/run-as-root.conf
# Workaround https://github.com/rfjakob/earlyoom/issues/270
[Service]
DynamicUser=false
