# etc/systemd/system/chrony.service.d/mullvad-exclude.conf
# Split-tunneling Chrony outside of Mullvad VPN (won't work without it)

[Unit]
# Don't start until mullvad-daemon is online to allow excluding
Wants=mullvad-daemon.service
After=network.target mullvad-daemon.service

[Service]
# Wait a minute before actually starting Chrony so Mullvad has surely had
# time to start and thus not fail at excluding (cannot set the cgroup, no
# such file or directory
ExecStartPre=/usr/bin/sleep 60
# Empty the main chrony.service ExecStart line so it can be overwritten
ExecStart=
ExecStart=/usr/bin/mullvad-exclude /usr/sbin/chronyd $DAEMON_OPTS
