# etc/systemd/system/chrony.service.d/mullvad-exclude.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mullvad-exclude /usr/sbin/chronyd $DAEMON_OPTS
