# etc/systemd/system/sailfish/ydns-simple.service
[Unit]
Description=ydns-simple service

[Service]
Type=simple
ExecStart=/usr/local/bin/ydns-simple

[Install]
WantedBy=ydns-simple.timer
