# etc/systemd/system/aminda-nocron-reboot.service.d/swapoff.conf
# Disables all swap. Very questionable.
[Service]
ExecStart=-/usr/sbin/swapoff -a
