# etc/systemd/system/ydns6/ydns6.service
[Unit]
Description=ydns6 service

[Service]
Type=simple
ExecStart=/home/nemo/.local/bin/ydns6

[Install]
WantedBy=ydns6.timer
