# etc/systemd/system/systemd-resolved.service.d/conflicts-avahi.conf
[Unit]
Conflicts=avahi-daemon.service avahi-daemon.socket
