# etc/systemd/system/i2pd.service.d/override.conf
[Service]
ExecStartPre=/bin/rm -f /var/lib/i2pd/routerinfo
ExecStartPre=/bin/sleep 15
