# etc/systemd/system/snapd.refresh.service.d/override.conf
[Unit]
Requires=network-online.target
After=network-online.target
