# etc/systemd/system/qbittorrent-nox@.service.d/webui-port.conf
[Service]
ExecStart=
# Contained in `systemd-path search-binaries-default`
ExecStart=qbittorrent-nox --webui-port=8081
