# etc/profile.d/99-maliit.sh
# Using maliit
export XMODIFIERS=@im="maliit"
export GTK_IM_MODULE="maliit"
export QT_IM_MODULE="maliit"
