# etc/profile.d/mikaela.sh
# Map compose to left Windows key and kill x with ctrl alt backspace and
# disable unbreakable space
setxkbmap -option compose:lwin -option terminate:ctrl_alt_bksp -option nbsp:none

# Steam won't work with it's included runtime
#export STEAM_RUNTIME=0
