# etc/systemd/system/miredo.service.d/override.conf
[Unit]
# If unbound is used for DNS, also add unbound.service
After=network-online.target
