# etc/systemd/system/yggdrasil.service.d/never-fail.conf
# To never stop retrying to start
[Service]
Restart=always
RestartSec=5s
