# etc/xdg/autostart/a-darkman.desktop
[Desktop Entry]
Type=Application
Exec=darkman run
Name=darkman
Icon=darkman
