# etc/xdg/autostart/a-telegram.desktop
[Desktop Entry]
Name=TELEGRAM
Exec=/usr/local/bin/telegram-desktop
Icon=telegram
Type=Application
