# etc/xdg/autostart/clipit.desktop
[Desktop Entry]
Name=ClipIt
Icon=clipit-trayicon
Exec=/usr/bin/clipit
Terminal=false
Type=Application
