# etc/xdg/autostart/com.wire.WireDesktop.desktop
[Desktop Entry]
Name=Wire
Exec=/usr/bin/flatpak run com.wire.WireDesktop --startup
Type=Application
Icon=com.wire.WireDesktop
