# etc/xdg/autostart/dino.desktop
[Desktop Entry]
Name=Dino
Exec=dino
Icon=dino
Type=Application
