# etc/xdg/autostart/im.riot.Riot.desktop
[Desktop Entry]
Name=Element
Exec=/usr/bin/flatpak run im.riot.Riot --hidden
Type=Application
Icon=im.riot.Riot
