# etc/xdg/autostart/im.riot.Riot.desktop
[Desktop Entry]
Name=Riot
Exec=/usr/bin/flatpak im.riot.Riot --hidden
Type=Application
Icon=im.riot.Riot
