# etc/xdg/autostart/jami.desktop
[Desktop Entry]
Name=Jami
Exec=jami-gnome -r
Icon=jami
Type=Application
