# etc/xdg/autostart/linphone.desktop
[Desktop Entry]
Type=Application
Exec=linphone --iconified
Icon=linphone
Hidden=false
Name=linphone
