# etc/xdg/autostart/mullvad-vpn.desktop
[Desktop Entry]
Name=Mullvad VPN
Exec="/opt/Mullvad VPN/mullvad-vpn" %U
Type=Application
Icon=mullvad-vpn
