# etc/xdg/autostart/Nextcloud.desktop
[Desktop Entry]
Name=Nextcloud
Exec=/usr/bin/nextcloud
Icon=nextcloud
Type=Application
