# etc/xdg/autostart/org.signal.Signal.desktop
[Desktop Entry]
Name=Signal Tray
Exec=/usr/bin/flatpak run org.signal.Signal --start-in-tray
Type=Application
Icon=org.signal.Signal
