# etc/xdg/autostart/org.gajim.Gajim.desktop
[Desktop Entry]
Name=Gajim
Exec=/usr/bin/flatpak run org.gajim.Gajim
Type=Application
Icon=org.gajim.Gajim
