# etc/xdg/autostart/telegramdesktop.desktop
[Desktop Entry]
Name=Telegram Desktop
Exec=telegram-desktop -many -startintray
Icon=telegram
Type=Application
