# etc/xdg/autostart/telegramdesktop.desktop
[Desktop Entry]
Name=Telegram Desktop
Exec=env LANG=en telegram-desktop -startintray
Icon=telegram
Type=Application
