# etc/systemd/system/aminda-duperemove.timer
[Unit]
Description=Run Aminda's deduplication tasks

[Timer]
OnBootSec=0
Persistent=true
RemainAfterElapse=false

[Install]
WantedBy=timers.target
