# etc/systemd/system/gitea-restart.service
[Unit]
Description=Gitea restarter

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart gitea.service --quiet
