# etc/systemd/system/matterbridge-restart.service
[Unit]
Description=Matterbridge restarter

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart matterbridge.service --quiet
