# etc/systemd/system/oidentd.socket
[Unit]
Description=Ident (RFC 1413) socket

[Socket]
ListenStream=113
BindIPv6Only=both
Accept=yes

[Install]
WantedBy=sockets.target
