# etc/systemd/system/syncplay-server.service
[Unit]
Description=Syncplay server
After=network.target

[Service]
ExecStart=/opt/syncplay/syncplayServer.py --port 1234 --motd-file /opt/syncplay/motd.txt --salt autogenerated-by-running-without-this-flag-CHANGEME --password NotInThisConfigFileNotSorry
Restart=always
User=syncplay

[Install]
WantedBy=default.target
