# etc/systemd/system/syncplay-server.service
[Unit]
Description=Syncplay server
After=network.target

[Service]
Environment=PATH=/opt/syncplay/venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/games:/usr/games
ExecStart=/opt/syncplay/venv/bin/python3 /opt/syncplay/syncplayServer.py --port 1234 --motd-file /opt/syncplay/motd.txt --salt autogenerated-by-running-without-this-flag-CHANGEME --password NotInThisConfigFileNotSorry --tls /opt/syncplay/ssl
Restart=always
User=syncplay

[Install]
WantedBy=default.target
