# etc/systemd/system/ydns-simple.timer
[Unit]
Description=Run ydns-simple every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
