# etc/systemd/system/zeronet.service
[Unit]
Description=ZeroNet
After=network.target

[Service]
ExecStart=/usr/bin/python2 /home/mikaela/src/github/HelloZeroNet/ZeroNet/zeronet.py
Restart=on-failure
User=mikaela

[Install]
WantedBy=default.target
